# src/gray386/etc/inittab
::sysinit:/etc/init.d/rc
::respawn:/bin/cttyhack /bin/sh
::restart:/sbin/init
::ctrlaltdel:/sbin/reboot
::shutdown:/bin/umount -a -r
